# supervisor/fancontrol.conf
[program:fancontrol]
#command=/usr/local/bin/fancontrol.py
command=/home/pi/scripts/system/fancontrol.sh
autostart=true
autorestart=true
startretries=3
